# s6/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf